# files/etc/profile.d/history_format.sh
# shellcheck shell=sh

HISTTIMEFORMAT="|%F %R| "

